# deploy/com.javis.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.javis.hub</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>-m</string>
        <string>src.voice_hub</string>
        <string>--server</string>
        <string>ws://192.168.219.106:8765</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8766</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/d43103/Workspace/projects/javis</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/d43103/Workspace/projects/javis</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/javis-hub.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/javis-hub-error.log</string>
</dict>
</plist>
