# scripts/ai.javis.mic-client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.javis.mic-client</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/d43103/Workspace/projects/javis/scripts/run_mic_client_service.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/d43103/Workspace/projects/javis</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>JAVIS_SERVER_URL</key>
    <string>ws://100.67.60.57:8765</string>
    <key>JAVIS_SESSION_ID</key>
    <string>mac-service</string>
    <key>JAVIS_CLIENT_LOG_FILE</key>
    <string>/Users/d43103/Workspace/projects/javis/logs/client-service.log</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/d43103/Workspace/projects/javis/logs/client-service.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/d43103/Workspace/projects/javis/logs/client-service.err.log</string>
</dict>
</plist>
